# kallaxsizer.service
[Unit]
Description=Kallax Sizer - Board Game Shelf Calculator
After=network.target

[Service]
Type=simple
User=onewing
WorkingDirectory=/opt/kallaxsizer
ExecStart=/opt/kallaxsizer/.venv/bin/python3 /opt/kallaxsizer/server.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
